# stage3/05-prep-revvy/files/revvy.service
[Unit]
Description=Revvy service
StartLimitIntervalSec=1
After=hciuart.service

[Service]
Type=simple
Restart=always
RestartSec=1
User=pi
ExecStart=/usr/bin/env python3 /home/pi/RevvyFramework/launch_revvy.py

[Install]
WantedBy=bluetooth.target
