# stage3/05-prep-revvy/files/revvy.service
[Unit]
Description=Revvy service
After=bluetooth.service network.target sound.target hciuart.service polkit.service
Wants=bluetooth.service network.target sound.target multi-user.target hciuart.service dhcpcd.service systemd-logind.service dbus.service
StartLimitIntervalSec=1

[Service]
Type=simple
Restart=always
RestartSec=1
User=pi
ExecStart=/usr/bin/env python3 /home/pi/RevvyFramework/launch_revvy.py

[Install]
WantedBy=bluetooth.target
